# srvMonitor2/etc/systemd/system/srvMonitor2.service
[Unit]
Description=srvMonitor2

[Service]
ExecStart=/usr/local/bin/srvMonitor2.sh

[Install]
WantedBy=multi-user.target
