# srvmonitor2/etc/systemd/system/srvmonitor2.service
[Unit]
Description=srvmonitor2

[Service]
ExecStart=/usr/bin/srvmonitor2

[Install]
WantedBy=multi-user.target
